# restic-restore@.service
[Unit]
Description=Restic Backup

[Service]
EnvironmentFile=/%I
Type=oneshot
ExecStart=/usr/bin/restic restore latest --target /
