# restic-restore@.service
[Unit]
Description=Restic Backup

[Service]
EnvironmentFile=/etc/restic/%I.env
Type=oneshot
ExecStart=/usr/bin/restic restore latest --target /
